# libgalileo2/libgalileo2_service.1
#! /bin/sh
### BEGIN INIT INFO
# Provides:          libgalileo2_service
# Required-Start:
# Should-Start:
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Configures the libgalileo2_service functionalities
# Description:       Configures the libgalileo2_service functionalities
### END INIT INFO

case "$1" in
    start|restart|force-reload)


		## Begin of IO3 (PWM1) ##
		# As a output SERVO

			if [ ! -d /sys/class/pwm/pwmchip0/pwm1 ] ; then
				echo -n "1" > /sys/class/pwm/pwmchip0/export
			fi
			chgrp pwm /sys/class/pwm/pwmchip0/device/pwm_period
			chmod g+w /sys/class/pwm/pwmchip0/device/pwm_period
			chgrp pwm /sys/class/pwm/pwmchip0/pwm1/duty_cycle
			chmod g+w /sys/class/pwm/pwmchip0/pwm1/duty_cycle
			chgrp pwm /sys/class/pwm/pwmchip0/pwm1/enable
			chmod g+w /sys/class/pwm/pwmchip0/pwm1/enable

			# gpio16 = 0 = out
			if [ ! -d /sys/class/gpio/gpio16 ] ; then
				echo -n "16" > /sys/class/gpio/export
			fi
			echo -n "out" > /sys/class/gpio/gpio16/direction
			echo -n "0" > /sys/class/gpio/gpio16/value

			# gpio17 = in = no pull-up nor pull-down
			if [ ! -d /sys/class/gpio/gpio17 ] ; then
				echo -n "17" > /sys/class/gpio/export
			fi
			echo -n "in" > /sys/class/gpio/gpio17/direction

			if [ ! -d /sys/class/gpio/gpio76 ] ; then
				echo -n "76" > /sys/class/gpio/export
			fi
			echo -n "0" > /sys/class/gpio/gpio76/value


			# gpio64 = 1 */
			if [ ! -d /sys/class/gpio/gpio64 ] ; then
				echo -n "64" > /sys/class/gpio/export
			fi
			echo -n "1" > /sys/class/gpio/gpio64/value


		## End of IO3 (PWM1) ##

		## Begin of IO3 ##
		# As a output LED

				#if [ ! -d /sys/class/gpio/gpio62 ] ; then
				#	echo -n "62" > /sys/class/gpio/export
				#fi
				#echo -n "out" > /sys/class/gpio/gpio62/direction


				#if [ ! -d /sys/class/gpio/gpio76 ] ; then
				#	echo -n "76" > /sys/class/gpio/export
				#fi
				#echo -n "0" > /sys/class/gpio/gpio76/value

				#if [ ! -d /sys/class/gpio/gpio64 ] ; then
				#	echo -n "64" > /sys/class/gpio/export
				#fi
				#echo -n "out" > /sys/class/gpio/gpio64/direction
				#echo -n "0" > /sys/class/gpio/gpio64/value

				#chgrp gpio /sys/class/gpio/gpio62/value
				#chmod g+rw /sys/class/gpio/gpio62/value


		## End of IO3 ##

		## Begin of IO6 ##
		# As a output LED

			#if [ ! -d /sys/class/gpio/gpio1 ] ; then
			#	echo -n "1" > /sys/class/gpio/export
			#fi
			#echo -n "out" > /sys/class/gpio/gpio1/direction
			#
			#
			#if [ ! -d /sys/class/gpio/gpio20 ] ; then
			#	echo -n "20" > /sys/class/gpio/export
			#fi
			#echo -n "out" > /sys/class/gpio/gpio20/direction
			#echo -n "0" > /sys/class/gpio/gpio20/value
			#
			#if [ ! -d /sys/class/gpio/gpio68 ] ; then
			#	echo -n "68" > /sys/class/gpio/export
			#fi
			#echo -n "0" > /sys/class/gpio/gpio68/value
			#
			#chgrp gpio /sys/class/gpio/gpio1/value
			##chmod g+rw /sys/class/gpio/gpio1/value


		## End of IO6 ##

		## Begin of IO6 ##
		# As a input button

			if [ ! -d /sys/class/gpio/gpio1 ] ; then
				echo -n "1" > /sys/class/gpio/export
			fi
			echo -n "in" > /sys/class/gpio/gpio1/direction

			if [ ! -d /sys/class/gpio/gpio20 ] ; then
				echo -n "20" > /sys/class/gpio/export
			fi
			echo -n "out" > /sys/class/gpio/gpio20/direction
			echo -n "1" > /sys/class/gpio/gpio20/value

			if [ ! -d /sys/class/gpio/gpio21 ] ; then
				echo -n "21" > /sys/class/gpio/export
			fi
			echo -n "out" > /sys/class/gpio/gpio21/direction
			echo -n "0" > /sys/class/gpio/gpio21/value

			if [ ! -d /sys/class/gpio/gpio68 ] ; then
				echo -n "68" > /sys/class/gpio/export
			fi
			echo -n "0" > /sys/class/gpio/gpio68/value

			chgrp gpio /sys/class/gpio/gpio1/value
			chmod g+rw /sys/class/gpio/gpio1/value
			chgrp gpio /sys/class/gpio/gpio1/edge
			chmod g+rw /sys/class/gpio/gpio1/edge

		## End of IO6 ##

		## Begin of AD 0 ## 

			# gpio49 = in = neither pull-up nor pull-down
			if [ ! -d /sys/class/gpio/gpio49 ] ; then
				echo -n "49" > /sys/class/gpio/export
			fi
			echo -n "in" > /sys/class/gpio/gpio49/direction
				
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage0_raw
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage0_raw
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage0_scale
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage0_scale
			
		## End of AD 0 ##

		## Begin of AD 1 ## 

			# gpio51 = in = neither pull-up nor pull-down
			if [ ! -d /sys/class/gpio/gpio51 ] ; then
				echo -n "51" > /sys/class/gpio/export
			fi
			echo -n "in" > /sys/class/gpio/gpio51/direction
				
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage1_raw
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage1_raw
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage1_scale
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage1_scale
			
		## End of AD 1 ##

		## Begin of AD 2 ## 

			# gpio53 = in = neither pull-up nor pull-down
			if [ ! -d /sys/class/gpio/gpio53 ] ; then
				echo -n "53" > /sys/class/gpio/export
			fi
			echo -n "in" > /sys/class/gpio/gpio53/direction
				
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage2_raw
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage2_raw
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage2_scale
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage2_scale
			
		## End of AD 2 ##

		## Begin of AD 3 ## 

			# gpio55 = in = neither pull-up nor pull-down
			if [ ! -d /sys/class/gpio/gpio55 ] ; then
				echo -n "55" > /sys/class/gpio/export
			fi
			echo -n "in" > /sys/class/gpio/gpio55/direction
				
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage3_raw
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage3_raw
			chgrp adc /sys/bus/iio/devices/iio:device0/in_voltage3_scale
			chmod g+r /sys/bus/iio/devices/iio:device0/in_voltage3_scale
			
		## End of AD 3 ##

		## Begin of EEPROM ##
			chgrp eeprom /sys/bus/i2c/devices/0-0054/eeprom
			chmod g+rw /sys/bus/i2c/devices/0-0054/eeprom
		## End of EEPROM ##

		## Begin of LCD ##
			# Pull-up on IO18
			if [ ! -d /sys/class/gpio/gpio57 ] ; then
					echo -n "57" > /sys/class/gpio/export
			fi
			echo "out" > /sys/class/gpio/gpio57/direction
			echo "1" > /sys/class/gpio/gpio57/value

			# Select I2C on IO18 and IO19
			if [ ! -d /sys/class/gpio/gpio60 ] ; then
					echo -n "60" > /sys/class/gpio/export
			fi
			echo "out" > /sys/class/gpio/gpio60/direction
			echo "0" > /sys/class/gpio/gpio60/value

			# Pull-up on IO19
			if [ ! -d /sys/class/gpio/gpio59 ] ; then
					echo -n "59" > /sys/class/gpio/export
			fi
			echo "out" > /sys/class/gpio/gpio59/direction
			echo "1" > /sys/class/gpio/gpio59/value

			chgrp i2c /dev/i2c-0
			chmod g+rw /dev/i2c-0
		## End of LCD ##

		
	;;
    stop)

		## Begin of IO3 (PWM1) ## 
			echo -n "0" > /sys/class/pwm/pwmchip0/pwm1/enable
			echo -n "76" > /sys/class/gpio/unexport
			echo -n "64" > /sys/class/gpio/unexport
			echo -n "17" > /sys/class/gpio/unexport
			echo -n "1" > /sys/class/gpio/gpio16/value
			echo -n "16" > /sys/class/gpio/unexport
			echo -n "1" > /sys/class/pwm/pwmchip0/unexport
		## End of IO3 ## 

		## Begin of IO3 (LED) ## 
			#echo -n "in" > /sys/class/gpio/gpio62/direction
			#echo -n "62" > /sys/class/gpio/unexport
			#echo -n "76" > /sys/class/gpio/unexport
			#echo -n "64" > /sys/class/gpio/unexport
		## End of IO3 ## 

		## Begin of IO6 (LED) ## 
			echo -n "in" > /sys/class/gpio/gpio1/direction
			echo -n "1" > /sys/class/gpio/unexport
			echo -n "20" > /sys/class/gpio/unexport
			echo -n "68" > /sys/class/gpio/unexport
		## End of IO6 ## 


		## Begin of AD 0 ## 
			echo -n "49" > /sys/class/gpio/unexport
		## End of AD 0 ## 

		## Begin of AD 1 ## 
			echo -n "51" > /sys/class/gpio/unexport
		## End of AD 1 ## 

		## Begin of AD 3 ## 
			echo -n "53" > /sys/class/gpio/unexport
		## End of AD 3 ## 

		## Begin of AD 4 ## 
			echo -n "55" > /sys/class/gpio/unexport
		## End of AD 4 ## 

		## Begin of LCD ##
			# No pull-up on IO18
			echo "in" > /sys/class/gpio/gpio57/direction
			echo "57" > /sys/class/gpio/unexport

			# GPIO on IO18 and IO19
			echo "1" > /sys/class/gpio/gpio60/value
			echo "60" > /sys/class/gpio/unexport

			# No pull-up on  IO19
			echo "in" > /sys/class/gpio/gpio59/direction
			echo "59" > /sys/class/gpio/unexport
		## End of LCD ##

		
	;;
    status)
		## Begin of IO3 (pwm1) ##
			echo -n Period:
			cat /sys/class/pwm/pwmchip0/device/pwm_period
			echo -n Duty Cicle:
			cat /sys/class/pwm/pwmchip0/pwm1/duty_cycle
			echo -n Enable:
			cat /sys/class/pwm/pwmchip0/pwm1/enable
		## End of IO3 (pwm1) ##
	;;
    *)
	echo -n "Usage: $0 "
	echo "{start|stop|restart|force-reload|status}"
	exit 1
	;;
esac

exit 0
